# deploy/vpn-pack-manager.service
[Unit]
Description=VPN Pack Manager (UniFi)
After=tailscaled.service ubios-udapi-server.service unifi-core.service
Wants=tailscaled.service

[Service]
ExecStart=/persistent/vpn-pack/bin/vpn-pack-manager \
    --socket=/run/tailscale/tailscaled.sock
Restart=on-failure
RestartSec=5
NoNewPrivileges=yes
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes

[Install]
WantedBy=multi-user.target
